# crontab
*/5 *	* * *	root	sh /opt/unbound/unboundSend.sh {IP} {ZBX DEV} > /dev/null
